# roles/updater/templates/ansible_pull.service.j2
# {{ ansible_managed }}
[Unit]
Description=Ansible Pull for {{ item.name }}
Wants={{ item.name }}.timer

[Service]
ExecStart={{ updater_ansible_venv }}/bin/ansible-pull \
    -i {{ ansible_hosts_file }} \
    -C {{ item.checkout }} \
    -U {{ item.url }} \
    {{ item.playbook }}

Slice=ansible_pull.slice

[Install]
WantedBy=multi-user.target