# roles/updater/templates/updater.timer.j2
# {{ ansible_managed }}
[Unit]
Description={{ item.name }} timer
Wants=network.target network-online.target
After=network-online.target time-set.target time-sync.target
[Timer]
Unit={{item.name}}.service
Persistent=true
{{ item.timer }}
RandomizedDelaySec=600
OnBootSec=1s
RandomizedDelaySec=60s
[Install]
WantedBy=timers.target