# roles/install/templates/systemd/mqtt.path.j2
[Unit]
Description="Watch {{ printnanny_license_file }}"

[Path]
PathExists={{ printnanny_license_file }}
Unit=printnanny-mqtt.service
MakeDirectory=true

[Install]
WantedBy=multi-user.target