# roles/install/templates/systemd/printnanny.target.j2
[Unit]
Description=PrintNanny Services Available
After=network.target network-online.target
Wants=network.target network-online.target

[Install]
WantedBy=multi-user.target