# roles/printnanny/templates/systemd/printnanny-update.target.j2
[Unit]
Description=Print Nanny Updated
After=network.target printnanny-metadata.target
Wants=network.target printnanny-metadata.target printnanny-update.service

[Install]
WantedBy=multi-user.target