# roles/printnanny/templates/systemd/printnanny.target.j2
[Unit]
Description=Print Nanny Services
After=network.target network-online.target
Wants=network.target network-online.target

[Install]
WantedBy=multi-user.target