# roles/printnanny/templates/systemd/printnanny-update.target.j2
[Unit]
Description=Print Nanny Updated
After=syslog.target network.target
Wants=syslog.target network.target printnanny-update.service

[Install]
WantedBy=multi-user.target